# config/mysqld.conf
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --socket=/var/lib/mysql/mysql.sock --pid-file=/var/run/mysqld/mysqld.pid
autostart=true
autorestart=true
